# udev-rules/99-new-display.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/qlio/.config/udev-rules/99-new-display.sh"
